# hermes/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m hermes_cli.main gateway run --replace
WorkingDirectory=/home/YOURUSER/.local/lib/python3.12/site-packages
Environment=PATH=/home/YOURUSER/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=VIRTUAL_ENV=/home/YOURUSER/system/.venv
Environment=HERMES_HOME=/home/YOURUSER/system/hermes
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
